# ali.service
[Unit]
Description=Automatic Backlight

[Service]
ExecStart=/usr/bin/ali

[Install]
WantedBy=multi-user.target
